# profile.d/cse110-27.sh
#!/usr/bin/bash.exe
echo "Load cse110-27 profile!"