# profile.d/installClasses-alias.sh
#!/usr/bin/bash.exe
alias icl=installClasses.sh