# profile.d/path.sh
#!/usr/bin/bash.exe
PATH=/c/Users/Public/OneDrive/Documents/Projects/BYU-I/cse110-27/scripts:${PATH}